# systemd/bk.service
[Unit]
Description=Backup

[Service]
Type=simple
ExecStart=/usr/bin/bk run /etc/bk.toml
User=root
StandardOutput=journal
StandardError=journal
